# icinga2.conf
[program:icinga2]
command=/opt/icinga2-command
