# icinga-web.conf
[program:icinga-web]
command=/opt/icinga-web-command
